# docker/allinone/supervisord.conf
# supervisord programs for the all-in-one DNSTAP2 container.
# All four processes log to the container's stdout/stderr (logfile_maxbytes=0
# disables rotation so `docker logs dnstap2` shows everything). priority sets
# rough start order: stores first, then receiver/scraper, then Grafana.

[program:loki]
command=/usr/local/bin/loki -config.file=/etc/loki/loki-config.yml
priority=10
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:prometheus]
command=/usr/local/bin/prometheus --config.file=/etc/prometheus/prometheus.yml --storage.tsdb.path=/var/lib/prometheus --web.listen-address=0.0.0.0:9090
priority=20
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:dnscollector]
command=/usr/local/bin/dnscollector -config /etc/dnscollector/config.yml
priority=20
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:grafana]
command=/usr/share/grafana/bin/grafana-server --homepath=/usr/share/grafana --config=/etc/grafana/custom.ini
priority=30
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
